# services/fcscript.service
[Unit]
Description=Controle Financeiro

[Service]
ExecStart=/usr/bin/fcscript.sh

[Install]
WantedBy=multi-user.target
